# config-global/NetworkManager/dispatcher.d/99-resolv.conf-head_and_tail
#!/bin/bash

resolvconf='/etc/resolv.conf'
cat "$resolvconf"{.head,,.tail} 2>/dev/null > "$resolvconf".tmp
mv -f "$resolvconf".tmp "$resolvconf"
